# 09bash/crontab
MAILTO=admin@test.ru

0   *   *   *   *    cd /dir && ./nginx-stat.sh
